# Master_Alarm/service/slt-master-alarm.service
# =============================================================================
# SLT Lift Master Alarm System — systemd Service
# =============================================================================
# This service runs the master alarm monitor continuously on boot.
#
# Installation:
#   sudo cp slt-master-alarm.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable slt-master-alarm
#   sudo systemctl start slt-master-alarm
#
# Management:
#   sudo systemctl status slt-master-alarm     # Check status
#   sudo systemctl stop slt-master-alarm       # Stop service
#   sudo systemctl restart slt-master-alarm    # Restart service
#   journalctl -u slt-master-alarm -f          # Follow live logs
# =============================================================================

[Unit]
Description=SLT Lift Master Alarm Monitor — Industrial Emergency System
Documentation=https://github.com/SLT/SLT_Emergency_Alarm
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=gpio

# Working directory and executable
WorkingDirectory=/opt/slt-master-alarm
ExecStart=/opt/slt-master-alarm/venv/bin/python -u src/main.py

# Environment variables
# IMPORTANT: Replace with your actual InfluxDB API token
Environment=INFLUXDB_TOKEN=your-influxdb-api-token-here
Environment=ALARM_CONFIG_PATH=/opt/slt-master-alarm/config/config.yaml

# Restart policy — always restart on failure for 24/7 operation
Restart=always
RestartSec=10
StartLimitIntervalSec=300
StartLimitBurst=5

# Unbuffered Python output for real-time journal logging
Environment=PYTHONUNBUFFERED=1

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=slt-master-alarm

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/var/log/slt-master-alarm
ProtectHome=true

# Resource limits
MemoryMax=128M
CPUQuota=25%

[Install]
WantedBy=multi-user.target
